# packer/config/setup_boot.service
[Unit]
Description=Execute user-datas to setup the machine
After=network-online.target
[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/bash /ops/config/bootscript.sh
[Install]
WantedBy=multi-user.target